# configs/etc/systemd/system/rclone-sync@.service
[Unit]
Description=Automated rclone sync for %I

[Service]
User=root
Type=oneshot
Environment="JOB=%I"
EnvironmentFile=/etc/rclone-sync/default.env
EnvironmentFile=/etc/rclone-sync/%I.env
ExecStart=/usr/local/bin/rclone-sync.sh

[Install]
WantedBy=multi-user.target