# configs/etc/systemd/system/rclone-sync@.timer
[Unit]
Description=daily rclone sync for %I

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target