# configs/etc/update-motd.d/50-quote
#!/bin/sh
echo
printf "\"It's a Magical World, Hobbes Ol' Buddy ... Let's Go Exploring!\"\n"